# service/xyz.openbmc_project.bmc-pairing-manager.service.in
[Unit]
Description=BMC Pairing Manager D-Bus Service

[Service]
ExecStart=/usr/bin/bmc-pairing-manager -c /etc/bmc-pairing-manager/bmc-pairing-manager.conf
Restart=always
Type=dbus
BusName=@DEFAULT_BUSNAME@
[Install]
WantedBy=network.target
